# twake-websockets/websockets.app.twake.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0770                 ; socket file mode (default 0700)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


[program:websockets.app.twake]
command                 = sh ./websocket_supervisor.sh
directory               = /etc/supervisor/
process_name            = websockets.app.twake
numprocs                = 1
autostart               = true
autorestart             = true
user                    = root
stdout_logfile          = /twake-core/var/logs/websockets.app.twake.out.log
stdout_logfile_maxbytes = 5MB
stderr_logfile          = /twake-core/var/logs/websockets.app.twake.err.log
stderr_logfile_maxbytes = 5MB
